# systemd/gr26-dash-update.timer
[Unit]
Description=Periodic gr26-dash auto-update check
Documentation=https://github.com/Gaucho-Racing/TCM-26

[Timer]
# Fire 10s after boot — the dash is already up by then (it's on
# graphical-session.target which beats this timer easily), and there's
# no harm trying eagerly: if the network isn't ready, the script
# fails fast (~3s connect timeout) and the timer retries on its
# regular cadence. After that, every 6 hours so cars left on for a
# long session still pick up new releases.
OnBootSec=10s
OnUnitActiveSec=6h
# If the system was off when a fire was scheduled, run it on next boot.
Persistent=true

[Install]
WantedBy=timers.target
