# systemd/gr26-dash-update.timer
[Unit]
Description=Periodic gr26-dash auto-update check
Documentation=https://github.com/Gaucho-Racing/TCM-26

[Timer]
# Fire 1 min after boot — gives the network and NTP time to actually
# come up. Trying at 10s was too eager; the update would always fail
# the first time on a freshly booted car. After that, every 6h so
# cars left on for a long session still pick up new releases.
OnBootSec=1min
OnUnitActiveSec=6h
# If the system was off when a fire was scheduled, run it on next boot.
Persistent=true

[Install]
WantedBy=timers.target
