# systemd/gr26-dash-update.service
[Unit]
Description=GR26 Dash auto-update check
Documentation=https://github.com/Gaucho-Racing/TCM-26
# Triggered by gr26-dash-update.timer 5 min after boot and every 6h
# after — never blocks boot or the dash launching. No
# network-online.target dep on purpose; if there's no internet we
# fail fast inside the script and the timer will retry next interval.

[Service]
Type=oneshot
ExecStart=/usr/local/bin/gr26-dash-update
TimeoutStartSec=180
StandardOutput=journal
StandardError=journal
